# coolify/cron/crontab
# Cron schedule for yt-agent maintenance hits.
# All times UTC (set by Dockerfile).
#
# DASHBOARD_INTERNAL_URL is the docker-network hostname of the Next.js
# service — we don't bounce through the public domain because there's
# no need to incur TLS + Caddy hop overhead for internal cron pokes.
#
# RENDER_TRIGGER_KEY is the shared secret the maintenance routes check
# (existing convention from the GitHub Actions workflows).

# Every UTC hour — queue jobs for channels whose run_at_hour matches
# the current hour. Channels with run_at_hour==null preserve the legacy
# behaviour and only fire at 09:00 UTC (default hour inside the route).
# %H = current UTC hour (busybox date). `\%` escapes the % that cron
# otherwise treats as a newline in the command line.
0 * * * *  curl -fsS -X POST -H "X-API-Key: $RENDER_TRIGGER_KEY" "$DASHBOARD_INTERNAL_URL/api/maintenance/scheduled-render?hour=$(date -u +\%-H)" >> /proc/1/fd/1 2>&1

# 05:00 UTC daily — clean up old runs/logs/errors per retention rules.
0 5 * * *  curl -fsS -X POST -H "X-API-Key: $RENDER_TRIGGER_KEY" "$DASHBOARD_INTERNAL_URL/api/maintenance/cleanup" >> /proc/1/fd/1 2>&1

# Every 3 min — probe queue depth AND wake Kaggle inline if needed.
# The ?wake=1 param makes the route POST to /api/backends/wake-kaggle
# whenever it decides needs_worker=true, so we don't wait up to 5 min
# for the GitHub Actions dispatch cron to catch up when a worker dies.
*/3 * * * *  curl -fsS -H "X-API-Key: $RENDER_TRIGGER_KEY" "$DASHBOARD_INTERNAL_URL/api/maintenance/needs-worker?wake=1" >> /proc/1/fd/1 2>&1

# Every 15 min — sweep stale backends rows (kaggle kernels that died
# mid-heartbeat, container restarts pre-stable-INSTANCE_ID) and orphan
# videos (bytes landed in MinIO but no runs_index row + older than
# retention). Prevents PB row bloat + MinIO storage creep.
# stale_seconds=600 → drop backends idle >10 min.
# orphan_retention=604800 → keep orphan videos for 7 days before deleting.
*/15 * * * *  curl -fsS -X POST -H "X-API-Key: $RENDER_TRIGGER_KEY" "$DASHBOARD_INTERNAL_URL/api/maintenance/cleanup-stale?stale_seconds=600&orphan_retention=604800" >> /proc/1/fd/1 2>&1

# 03:30 UTC daily — reclaim Docker build cache + dangling images on
# the host. Frequent image rebuilds (dashboard + side-worker) leave
# multi-GB BuildKit layer cache; without this it accumulates 15-20 GB
# per week and eats the VPS disk. Confirmed live 2026-07-10: freed
# 18 GB in a single prune after a busy dev day.
# Requires docker-cli in this image + /var/run/docker.sock mount from
# the host (see coolify/docker-compose.yml cron service).
30 3 * * *  (docker builder prune -a -f && docker image prune -f && docker container prune -f) >> /proc/1/fd/1 2>&1

# Every 2 hours — sweep leaked render work_dirs from the side-worker
# container's writable layer. housekeeping.finalize_run is supposed to
# rmtree these post-publish, but any bug in that path (a non-published
# run, a failed R2 mirror, a crash before finalize) leaves ~1.2 GB per
# run in /app/output/videos inside the container. This is the belt-and-
# suspenders — any dir older than 2h (well past a 5-min render) that's
# still there is a leak, safe to nuke. Confirmed live 2026-07-13:
# accumulated 31 GB across 24 unpublished-but-successful runs before
# this sweep was added (bug in finalize_run's `published` guard).
# 2026-07-17 (audit): the sweep now EXCLUDES the active run's dir. A
# long render (20-shot + Pollinations retries on 1.5 vCPU) can exceed
# 120 min, and the run dir's mtime stops advancing once its top-level
# entries exist — the old sweep could rm -rf an in-progress render's
# inputs mid-encode. The active run_id comes from data/run_state.json.
0 */2 * * *  SW=$(docker ps --format '{{.Names}}' | grep '^side-worker' | head -1); [ -n "$SW" ] && docker exec $SW sh -c 'ACTIVE=$(python3 -c "import json;d=json.load(open(\"data/run_state.json\"));print(d.get(\"run_id\",\"\") if d.get(\"status\") in (\"running\",\"starting\") else \"\")" 2>/dev/null); for d in $(find /app/output/videos -mindepth 1 -maxdepth 1 -type d -mmin +120 2>/dev/null); do [ "$(basename $d)" = "$ACTIVE" ] && continue; rm -rf "$d"; done' >> /proc/1/fd/1 2>&1

# Every 2 min — CONTAINER REVIVAL. The single biggest cause of the
# recurring "server timeout" was NOT proxy staleness: on 2026-07-20 the
# dashboard container was found Exited(143) for 2 hours with
# RestartPolicy=no, so nothing ever brought it back and every heal cycle
# had no upstream to route to. Restart policies are now set at create
# time, but this is the belt-and-braces: any expected app container
# found in a non-running state gets started immediately. Logs only when
# it actually revives something.
*/2 * * * *  for C in $(docker ps -a --filter "name=dashboard-" --filter "name=side-worker" --filter "name=pocketbase-" --filter "name=caddy-" --filter "name=minio-" --format '{{.Names}}'); do S=$(docker inspect -f '{{.State.Status}}' "$C" 2>/dev/null); if [ "$S" != "running" ]; then echo "container-revive: $C was $S — starting"; docker start "$C" >/dev/null 2>&1 && echo "container-revive: $C started"; fi; done >> /proc/1/fd/1 2>&1

# Every 3 min — gateway auto-heal. The recurring "Gateway Timeout" /
# HTTP 000 comes from Caddy's stale connection pool after ANY container
# rotation (deploy scripts bounce Caddy, but the gateway sometimes goes
# stale again minutes later — root cause is Caddy re-caching the old
# container IP via docker DNS before the swap fully settles). This
# watchdog probes the public login page from INSIDE the cluster via
# Traefik's published port; if it isn't HTTP 200/3xx, it runs the full
# heal cycle (Caddy restart + Traefik restart). Self-quieting: healthy
# probes log nothing. Confirmed necessary 2026-07-15 — third stale
# gateway in one session despite post-deploy bounces.
# Softened 2026-07-17 (audit): (a) requires TWO consecutive failed
# probes (state file) so a transient blip / in-progress deploy doesn't
# trigger a restart; (b) only restarts CADDY (app-scoped) on the first
# heal — coolify-proxy (SHARED Traefik fronting every app on the host)
# is only bounced if Caddy alone didn't fix it by the NEXT tick (third
# consecutive failure).
# 2026-07-18: probe now runs in a HOST-NETWORK container. The old probe
# ran from this cron container's bridge network, whose hairpin NAT path
# kept working while the eth0-side DNAT was wedged — so the watchdog
# saw 200 while every real external client got connection timeouts
# (confirmed live: external 000, internal 200, Traefik logs silent =
# connections never reached the application layer). Host-network curl
# traverses the same iptables DOCKER chain externals use.
# REVERTED to in-container curl (2026-07-18): the host-network probe
# variant reliably returned 000 on Oracle (NAT hairpin quirk from a
# host-netns container), which made the watchdog restart Caddy/Traefik
# every 6-9 minutes — a self-inflicted outage storm. The in-container
# probe can miss one rare external-only failure mode, but a blind spot
# beats an active saboteur. --max-time 15 for TLS headroom.
*/3 * * * *  CODE=$(curl -s -o /dev/null -w '%{http_code}' --max-time 15 --resolve yt-agent.thyker.online:443:132.226.86.61 https://yt-agent.thyker.online/login 2>/dev/null || echo 000); case "$CODE" in 2*|3*) rm -f /tmp/gwfail ;; *) N=$(($(cat /tmp/gwfail 2>/dev/null || echo 0) + 1)); echo "$N" > /tmp/gwfail; echo "gateway-heal: probe returned $CODE (consecutive fail #$N)"; if [ "$N" -eq 2 ]; then CADDY=$(docker ps --format '{{.Names}}' | grep '^caddy' | head -1); [ -n "$CADDY" ] && docker restart $CADDY && echo "gateway-heal: restarted caddy"; elif [ "$N" -ge 3 ]; then docker restart coolify-proxy && echo "gateway-heal: escalated — restarted coolify-proxy"; rm -f /tmp/gwfail; fi; esac >> /proc/1/fd/1 2>&1

# Trailing newline is required — Linux crontab silently ignores the last
# line if it doesn't end with one.
